# shadow_server/deploy/xinyu-shadow.service
[Unit]
Description=Xinyu Shadow AI Backend
After=network.target mysql.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/xinyu/shadow_server
Environment=PATH=/opt/xinyu/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/opt/xinyu/venv/bin/python -m uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
