# shadow_server/deploy/xinyu-flask.service
[Unit]
Description=Xinyu Flask API (User/Emotion/Tarot)
After=network.target mysql.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/xinyu/xinyu/server
Environment=MYSQL_HOST=43.143.169.226
Environment=MYSQL_PORT=3306
Environment=MYSQL_USER=root
Environment=MYSQL_PASSWORD=010824Wy?
Environment=MYSQL_DATABASE=xintujie
Environment=PATH=/opt/xinyu/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/opt/xinyu/venv/bin/python flask_api.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
